# pkg/pkg/pkg-plist
/usr/local/lib/peerbanhelper/PeerBanHelper.jar
/usr/local/etc/rc.d/peerbanhelper
